# filedforBahnprojekt/deploy/onpoint-backend.service.example
[Unit]
Description=OnPoint Flask backend
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/BahnProjekt/filedforBahnprojekt
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/BahnProjekt/.venv/bin/gunicorn --workers 1 --bind 127.0.0.1:5000 --timeout 120 backend:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
